# deploy/mitre-agent.service
[Unit]
Description=Universal MITRE AI Agent — control panel
Documentation=file:/opt/mitre-agent/docs/Architecture.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mitre
Group=mitre
WorkingDirectory=/opt/mitre-agent
EnvironmentFile=/etc/mitre-agent/agent.env
ExecStart=/opt/mitre-agent/venv/bin/python /opt/mitre-agent/app.py
Restart=on-failure
RestartSec=5

# --- hardening (read-only source access; the agent never deploys rules) ---
# StateDirectory guarantees /var/lib/mitre-agent exists (mode 0750, owned by
# the service user) on every start and is implicitly read-write.
StateDirectory=mitre-agent
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mitre-agent
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service
LockPersonality=true

[Install]
WantedBy=multi-user.target
